# config/systemd/user/thunar.service.d/custom.conf
[Unit]
PartOf=graphical-session.target

[Install]
WantedBy=xsession.target
